# deploy/mrc-refresh-mcp.service
[Unit]
Description=MRC Refresh MCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mrc
Group=mrc
WorkingDirectory=/opt/mrc-refresh-mcp
Environment=PATH=/opt/mrc-refresh-mcp/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=TOKENS_PATH=/etc/mrc-refresh-mcp/tokens.yml
Environment=AUDIT_LOG=/home/mrc/mrc-proxy/logs/contractor-access.log
Environment=MRC_REFRESH_BIN=/usr/local/bin/mrc-refresh
ExecStart=/opt/mrc-refresh-mcp/venv/bin/uvicorn server:app --host 127.0.0.1 --port 8765 --timeout-keep-alive 120
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/mrc/mrc-proxy/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
